# deployment/entrypoint/raw_nginx_cache/raw_nginx_cache_backend.supervisord.conf
[program:raw_nginx_cache_backend]
command=python /home/raw_nginx_cache/app.py
stdout_logfile = syslog
redirect_stderr = true
autorestart = true