# deployment/entrypoint/supervisord.conf
[unix_http_server]
file=/tmp/seekscale-entrypoint.supervisor.sock     ; path to your socket file
 
[supervisord]
logfile=/var/log/seekscale-entrypoint/supervisord.log ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/seekscale-entrypoint.pid          ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/seekscale-entrypoint/         ; where child log files will live
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
 
[supervisorctl]
serverurl=unix:///tmp/seekscale-entrypoint.supervisor.sock ; use a unix:// URL  for a unix socket

[program:raw_nginx_cache_backend]
command=/usr/local/share/seekscale/raw_nginx_cache/venv/bin/python /usr/local/share/seekscale/raw_nginx_cache/app.py
stdout_logfile = /var/log/seekscale-entrypoint/raw_nginx_cache_backend.log
redirect_stderr = true
autorestart = true

[program:smbproxy]
command=/usr/local/share/seekscale/smbproxy/venv/bin/python /usr/local/share/seekscale/smbproxy/__main__.py --shares-root /home/data/smbshares/WIN-DQ8IVRQPK51 --metadata-proxy-address 127.0.0.1 --fileserver-address gateway.seekscale.com --fileserver-port 61100
stdout_logfile = /var/log/seekscale-entrypoint/smbproxy.log
redirect_stderr = true
autorestart = true

[program:metadata_proxy]
command=/usr/local/share/seekscale/smbproxy/venv/bin/python /usr/local/share/seekscale/smbproxy/metadata_proxy/metadata_proxy.py
stdout_logfile = /var/log/seekscale-entrypoint/metadata_proxy.log
redirect_stderr = true
autorestart = true

[program:stunnel]
command=stunnel4 /etc/seekscale/stunnel.conf
stdout_logfile = /var/log/seekscale-entrypoint/stunnel.log
redirect_stderr = true
autorestart = true
