# deployment/gateway/supervisord.conf
[unix_http_server]
file=/tmp/seekscale-gateway.supervisor.sock     ; path to your socket file
 
[supervisord]
logfile=/var/log/seekscale-gateway/supervisord.log ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/seekscale-gateway.pid          ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/seekscale-gateway/         ; where child log files will live
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
 
[supervisorctl]
serverurl=unix:///tmp/seekscale-gateway.supervisor.sock ; use a unix:// URL  for a unix socket
 
[program:fileserver4-uploader]
command=/usr/local/share/seekscale/venv/bin/python /usr/local/share/seekscale/fileserver4_uploader.py 6107%(process_num)s
process_name=fileserver4-uploader-%(process_num)s
numprocs=5
stdout_logfile = /var/log/seekscale-gateway/fileserver4-6107%(process_num)s
redirect_stderr = true

[program:fileserver4-downloader]
command=/usr/local/share/seekscale/venv/bin/python /usr/local/share/seekscale/fileserver4_downloader.py
process_name=fileserver4-downloader-%(process_num)s
numprocs=5
stdout_logfile = /var/log/seekscale-gateway/fileserver4-downloader-%(process_num)s
redirect_stderr = true

[program:fileserver4-metadata]
command=/usr/local/share/seekscale/venv/bin/python /usr/local/share/seekscale/fileserver4_metadata.py 610%(process_num)02d
process_name=fileserver4-metadata-%(process_num)s
numprocs=30
stdout_logfile = /var/log/seekscale-gateway/fileserver4-610%(process_num)02d
redirect_stderr = true

#[program:mount_drives]
#command=/usr/local/bin/seekscale_gateway config mount-drives
#stdout_logfile = /var/log/seekscale-gateway/mount-drives.log
#redirect_stderr = true

[program:stunnel]
command=stunnel4 /etc/seekscale/stunnel.conf
stdout_logfile = /var/log/seekscale-gateway/stunnel.log
redirect_stderr = true
autorestart = true
